# docs/systemd/ibkr-watchdog@.service
# docs/systemd/ibkr-watchdog@.service
[Unit]
Description=IBKR API Watchdog for user %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/vega
Environment=PYTHONUNBUFFERED=1
Environment=WATCHDOG_INTERVAL_SEC=600
Environment=WATCHDOG_FAIL_THRESHOLD=3
ExecStart=/usr/bin/python3 scripts/watchdog_ibkr.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
