# ibkr-bridge/ibkr-bridge.service
[Unit]
Description=IBKR FastAPI Bridge (LIVE 7496)
After=network.target

[Service]
User=root
WorkingDirectory=/root/ibkr-bridge
Environment=PATH=/root/ibkr-bridge/.venv/bin
Environment=IB_HOST=127.0.0.1
Environment=IB_PORT=7496
Environment=IB_CLIENT_ID=9
Environment=IB_PAPER_ONLY=false
Environment=IBKR_MARKET_DATA_TYPE=1
Environment=BRIDGE_API_KEY=VegaTrading2025X
ExecStart=/root/ibkr-bridge/.venv/bin/uvicorn bridge:app --host 0.0.0.0 --port 8088
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
